# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
user=root
nodaemon=true
logfile=/dev/fd/1
logfile_maxbytes=0

[program:earthworm]
directory=/eworm/run/params
command=/eworm/earthworm_7.9/bin/startstop 
user=aqms
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
